# com.slowletter.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.slowletter.daily-update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/slowclaw/slowletter-pipeline/daily_update.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/slowclaw/slowletter-pipeline</string>

  <!-- 매일 08:30 실행 (원하면 Hour/Minute 수정) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <!-- git push까지 자동화하려면: daily_update.sh 실행 시 GIT_PUSH=1로 설정하거나, 여기에서 EnvironmentVariables로 넣을 수 있습니다 -->
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_daily_update.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_daily_update.err</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
