# com.slowletter.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.slowletter.daily-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/slowclaw/slowletter-pipeline/daily_update.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/slowclaw/slowletter-pipeline</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
