# com.slowletter.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.slowletter.daily-update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/slowclaw/slowletter-pipeline/daily_update_rag.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/slowclaw/slowletter-pipeline</string>

  <!-- 매일 08:30 실행 (원하면 Hour/Minute 수정) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>GIT_PUSH</key>
    <string>1</string>
    <!-- git push 포함. 필요하면 0으로 바꾸면 됩니다 -->
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_daily_update.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/slowclaw/slowletter-pipeline/data/logs/launchd_daily_update.err</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
